# moriod/DEBIAN/postinst
#!/bin/bash
# Post-Install script for moriod

set -e

MORIO_VERSION=__MORIO_VERSION__
MORIO_CONTAINER_TAG=__MORIO_CONTAINER_TAG__
MORIO_RELEASE_CHANNEL=__MORIO_RELEASE_CHANNEL__

# Force the correct version
echo $MORIO_VERSION > /etc/morio/moriod/version
tee /etc/morio/moriod/channel.env > /dev/null <<EOF
#
# This file is auto-generated by the moriod software package.
# Under normal circumstances, you should not edit it.
#
# It is installed/provided by the moriod package
# and will be updated when you update the package.
#
MORIO_CONTAINER_TAG=${MORIO_CONTAINER_TAG}
MORIO_RELEASE_CHANNEL=${MORIO_RELEASE_CHANNEL}
EOF


# Reload systemd
systemctl daemon-reload || true

# Enable services
systemctl enable docker || true
systemctl enable moriod || true

# Start services
systemctl start docker || true
systemctl start moriod || true

